# deploy/itam-backup.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
30 2 * * * root /opt/itam/app/deploy/backup-postgres.sh >> /var/log/itam-backup.log 2>&1
